# jarvis.service
[Unit]
Description=Jarvis AI Assistant
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'source %h/jarvis_env/bin/activate && python %h/jarvis/jarvis.py --text'
# ↑ use --text se não quiser que ele "fale" no servidor;
#   remova --text se quiser saída de voz

WorkingDirectory=%h/jarvis
EnvironmentFile=%h/jarvis/.env
Restart=on-failure
RestartSec=10
User=%u

[Install]
WantedBy=multi-user.target
